# oracle/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:listener]
command=sh /startuporacle.sh
startsecs=20
startretries=1

