# oracle/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:listener]
command=/bin/sh -c "/startuporacle.sh"


