# jenkins/supervisord.conf
[supervisord]
nodaemon=true

[program:ssh]
command=/usr/sbin/sshd -D

[program:startup-jenkins]
command= /bin/sh -c \"java -jar /opt/jenkins/jenkins.war\""

