# jfrog/supervisord.conf
[supervisord]
nodaemon=true

[program:ssh]
command=/usr/sbin/sshd -D

[program:jenkins]
command=service artifactory start
startretries=0
