# endeca/endeca312/supervisord.conf
[supervisord]
nodaemon=true

[program:ssh]
command=/usr/sbin/sshd -D

[program:platformservice|toolsandframeworks|cas]
command=su - docker -c "/bin/sh -c \"/startupendeca.sh\""